# upstart/5sur.conf
# /etc/init/goserver.conf
description "5sur server"
start on runlevel [2345]
stop on runlevel [06]
respawn

chdir /5sur/server

exec /5sur/server/5sur